# monitor/40-monitor-plug.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/home/$USER/autorandr.sh"
